# scripts/cron.d/backup
@daily root /home/ckan/src/ckanext-danepubliczne/scripts/backup.sh
