# scripts/cron.d/ckan-tracking
@hourly ckan workon ckan && paster --plugin=ckan tracking update
