# scripts/cron.d/ckan-tracking
@hourly ckan /home/ckan/src/ckanext-danepubliczne/scripts/ckan-tracking.sh >> /var/log/ckan/tracking.log 2>&1
