# scripts/cron.d/reports_qa
@daily ckan /home/ckan/.virtualenvs/ckan/bin/paster --plugin=ckanext-report report generate